# sv/mpdris/run
#!/bin/sh
mpdris --host /run/user/1000/mpd.sock
